# ipabakd.service
[Unit]
Description=A daemon for automatic freeipa backups

[Service]
ExecStart=/usr/sbin/ipabakd.sh 4

[Install]
WantedBy=multi-user.target
